# debian/postrm
#!/bin/sh -e

if [ "$1" = "remove" ] || if [ "$1" = "purge" ]
then
  update-rc.d nat6 defaults
fi
