# debian/postinst
#!/bin/sh -e

if [ "$1" = "configure" ]
then
  update-rc.d nat6 defaults
fi
