# systemd/squid-block.service.d/10ipv6-support-dependency.conf
[Unit]
PartOf=squid.service
ReloadPropagatedFrom=squid.service
Before=squid.service

