# systemd/squid-block.service.d/10ipv6-support-dependency.conf
[Unit]
PartOf=squid.service
StartPropagatedFrom=squid.service
ReloadPropagatedFrom=squid.service
RestartPropagatedFrom=squid.service
StopPropagatedFrom=squid.service
Before=squid.service

