# systemd/squid.service.d/10ipv6-support-dependency.conf
[Unit]
Wants=squid-block.service
Before=squid-block.service
